# backup-rsync-fbs.service
[Unit]
Description=A job to test the systemd scheduler

[Service]
Type=simple
ExecStart=backup.sh

[Install]
WantedBy=default.target
